# systemd/very-simple-mcp.user.service
[Unit]
Description=Very Simple MCP Server (User)

[Service]
ExecStart=%h/.npm-global/bin/very-simple-mcp
Restart=always

[Install]
WantedBy=default.target

